# vaultwarden/rootfs/etc/services.d/vaultwarden/run
#!/usr/bin/with-contenv bashio

# ===== INFO =====
bashio::log.info "-----------------------------------------"
bashio::log.info "Starting Vaultwarden add-on v1.0.5 (debug mode)"
bashio::log.info "-----------------------------------------"

# ===== Проверка утилит =====
bashio::log.info "Проверка наличия openssl и file..."
if ! command -v openssl >/dev/null 2>&1; then
    bashio::log.warning "openssl не найден! Некоторые функции будут сгенерированы без него."
else
    bashio::log.info "openssl найден."
fi

if ! command -v file >/dev/null 2>&1; then
    bashio::log.warning "file не найден! Не смогу показать архитектуру бинаря."
else
    bashio::log.info "file найден."
fi

# ===== Проверка бинаря Vaultwarden =====
bashio::log.info "Проверка бинаря /usr/bin/vaultwarden..."
if [[ ! -f "/usr/bin/vaultwarden" ]]; then
    bashio::log.error "/usr/bin/vaultwarden не найден!"
    ls -l /usr/bin || true
    exit 1
fi

if [[ ! -x "/usr/bin/vaultwarden" ]]; then
    bashio::log.error "/usr/bin/vaultwarden не исполняемый!"
    ls -l /usr/bin/vaultwarden
    exit 1
fi

if command -v file >/dev/null 2>&1; then
    bashio::log.info "Архитектура бинаря:"
    file /usr/bin/vaultwarden
fi

# ===== Отладка web-vault =====
bashio::log.info "=== DEBUG WEB-VAULT ==="
bashio::log.info "Текущая директория (PWD): $(pwd)"
bashio::log.info "Содержимое /web-vault:"
ls -la /web-vault || bashio::log.error "Папка /web-vault не существует!"
if [[ -f "/web-vault/index.html" ]]; then
    bashio::log.info "index.html найден — UI должен работать."
    bashio::log.info "Размер папки: $(du -sh /web-vault)"
else
    bashio::log.error "index.html НЕ найден — UI сломан!"
fi
bashio::log.info "========================"

# ===== Чтение конфигурации =====
DOMAIN=$(bashio::config 'domain' 'http://192.168.0.133:80')
ROCKET_PORT=$(bashio::config 'rocket_port' '80')
ADMIN_TOKEN=$(bashio::config 'admin_token')
DATA_FOLDER="/data"
DATABASE_URL="/data/db.sqlite3"
WEBSOCKET_ENABLED=$(bashio::config 'websocket_enabled' 'true')
SIGNUPS_ALLOWED=$(bashio::config 'signups_allowed' 'false')
INVITATIONS_ALLOWED=$(bashio::config 'invitations_allowed' 'true')
EMERGENCY_ACCESS_ALLOWED=$(bashio::config 'emergency_access_allowed' 'true')
SENDS_ALLOWED=$(bashio::config 'sends_allowed' 'true')
WEB_VAULT_FOLDER="/web-vault"
WEB_VAULT_ENABLED="true"

# ===== Генерация ADMIN_TOKEN =====
if [[ -z "$ADMIN_TOKEN" ]]; then
    if command -v openssl >/dev/null 2>&1; then
        ADMIN_TOKEN=$(openssl rand -base64 48) || true
        bashio::log.warning "ADMIN_TOKEN не задан — сгенерирован: ${ADMIN_TOKEN:0:20}..."
    else
        ADMIN_TOKEN="changeme"
        bashio::log.warning "openssl не найден — ADMIN_TOKEN установлен как 'changeme'"
    fi
fi

# ===== Экспорт переменных =====
export DOMAIN
export SENDS_ALLOWED
export EMERGENCY_ACCESS_ALLOWED
export INVITATIONS_ALLOWED
export SIGNUPS_ALLOWED
export WEBSOCKET_ENABLED
export DATABASE_URL
export DATA_FOLDER
export ADMIN_TOKEN
export ROCKET_PORT
export ROCKET_ADDRESS="0.0.0.0"
export WEB_VAULT_FOLDER
export WEB_VAULT_ENABLED

# ===== Вывод для отладки =====
bashio::log.info "=== DEBUG ENV EXPORTED ==="
bashio::log.info "DOMAIN: $DOMAIN"
bashio::log.info "ROCKET_PORT: $ROCKET_PORT"
bashio::log.info "DATA_FOLDER: $DATA_FOLDER"
bashio::log.info "DATABASE_URL: $DATABASE_URL"
bashio::log.info "ROCKET_ADDRESS: $ROCKET_ADDRESS"
bashio::log.info "WEBSOCKET_ENABLED: $WEBSOCKET_ENABLED"
bashio::log.info "SIGNUPS_ALLOWED: $SIGNUPS_ALLOWED"
bashio::log.info "INVITATIONS_ALLOWED: $INVITATIONS_ALLOWED"
bashio::log.info "EMERGENCY_ACCESS_ALLOWED: $EMERGENCY_ACCESS_ALLOWED"
bashio::log.info "SENDS_ALLOWED: $SENDS_ALLOWED"
bashio::log.info "ADMIN_TOKEN: ${ADMIN_TOKEN:0:20}..."
bashio::log.info "==========================="

# ===== Проверка доступности Vaultwarden =====
bashio::log.info "Проверка доступа к бинарю Vaultwarden..."
ls -l /usr/bin/vaultwarden || true

bashio::log.info "Попытка выполнить /usr/bin/vaultwarden --version..."
if ! /usr/bin/vaultwarden --version; then
    bashio::log.error "Выполнение --version завершилось с ошибкой. Возможно, неправильная архитектура или отсутствуют зависимости!"
fi

# ===== Запуск Vaultwarden =====
bashio::log.info "Исполняем /usr/bin/vaultwarden..."
exec /usr/bin/vaultwarden

